# scripts/guest-image/mkosi/mkosi.extra/etc/systemd/system/omnifs-frontend.service
[Unit]
Description=omnifs FUSE frontend runner
Requires=omnifs-seed-mount.service
After=omnifs-seed-mount.service systemd-modules-load.service

[Service]
Type=exec
# The seed conf is EnvironmentFile-shaped (KEY=VALUE lines): OMNIFS_ATTACH_ADDR,
# OMNIFS_ATTACH_TOKEN, and OMNIFS_READY_VSOCK_PORT (crates/omnifs-api/src/lib.rs).
# The runner dials host vsock on that port and writes `ready\n` once the FUSE
# mount is serving (crates/omnifs-namespace-wire/src/beacon.rs), which is how
# the krunkit backend's `mount_ready` observes guest readiness. No leading '-':
# a missing or unreadable file fails this unit's start with a clear journal
# message instead of running the frontend with no attach target.
EnvironmentFile=/mnt/seed/omnifs-seed.conf
ExecStart=/usr/local/bin/omnifs-fuse --mount-point /omnifs
StandardOutput=journal+console
StandardError=journal+console

[Install]
WantedBy=multi-user.target
